# config/crontab
0 0 * * * arkmanager update --warn --update-mods
0 0 * * * arkmanager backup